# overlay/etc/cron.d/elgg
# Taken from https://learn.elgg.org/en/stable/admin/cron.html
# Leverages TurnKey provided elgg-cli wrapper

* * * * * /usr/local/bin/turnkey-elgg-cli cron -q
